# rider-camera.service
[Unit]
Description=XGO Rider Camera (CSI camera -> MQTT image responder)
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/xgovenv/bin/python /home/pi/rider_camera.py
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
