# pi/rider-joystick.service
[Unit]
Description=XGO Rider DS4 controller -> MQTT (rider/control/line)
After=network.target mosquitto.service bluetooth.service rider-bridge.service
Wants=mosquitto.service rider-bridge.service

[Service]
Type=simple
User=pi
Environment=SDL_VIDEODRIVER=dummy
Environment=SDL_JOYSTICK_ALLOW_BACKGROUND_EVENTS=1
ExecStart=/home/pi/xgovenv/bin/python /home/pi/rider_controller.py
Restart=on-failure
RestartSec=3
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
